# docker/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock  ; the path to the socket file

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; main log file; default $CWD/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL for a unix socket

[program:octane]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan octane:start --port=8080
user=nobody
autostart=true
autorestart=true
startretries=9999
stdout_logfile=/var/www/html/storage/octane-stdout.log
stderr_logfile=/var/www/html/storage/octane-stderr.log

[program:supercronic]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/supercronic /etc/crontabs/nobody
autostart=true
autorestart=true
startretries=9999
user=nobody
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/supercronic.log
